# etc/systemd/system/hosting-autoscale.service
[Unit]
Description=Auto-scale Tor hosting instances based on account count

[Service]
Type=oneshot
ExecStart=/usr/bin/php8.5 /var/www/autoscale.php
