# etc/systemd/system/hosting-autoscale.timer
[Unit]
Description=Check and auto-scale hosting instances every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
